# deploy/systemd/apk-builder-retention.timer
[Unit]
Description=Run APK Builder workspace retention daily

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
